# supervisor_conf/bt_service.conf
[program:bt_service]
command=./service.py
directory=/root/eecs399/bluetooth_send/
environment=HOME="/root/"
priority=500
autostart=true
autorestart=true
