# supervisor_conf/OPUSH_port.conf
[program:OPUSH_port]
command=sdptool add --channel=12 OPUSH
priority=300
autostart=true
autorestart=false
startretries=1
