# supervisor_conf/SP_port.conf
[program:SP_port]
command=sdptool add --channel=1 SP
priority=200
autostart=true
autorestart=false
startretries=1
